# container-files/etc/services.d/sendmail/run
#!/usr/bin/wit-contenv /bin/sh

# Do the things needed (if any other than already set in configuration files) to 
# start the sendmail service.

/usr/sbin/sendmail -bD